# scripts/systemd/credentials-fetcher.service
[Unit]
Description=credentials-fetcher systemd service unit file.

[Service]
ExecStartPre=mkdir -p /var/credentials-fetcher/krbdir /var/credentials-fetcher/socket /var/credentials-fetcher/logging
ExecStart=/usr/sbin/credentials-fetcherd
Environment="CREDENTIALS_FETCHERD_STARTED_BY_SYSTEMD=1"
Environment="CF_CRED_FILE=/var/credentials-fetcher/credspec.json"
Type=notify
NotifyAccess=main
WatchdogSec=5s
Restart=on-failure

[Install]
WantedBy=multi-user.target
