# credentials-fetcher/scripts/systemd/credentials-fetcher.service
[Unit]
Description=credentials-fetcher systemd service unit file.

[Service]
ExecStart=/usr/sbin/credentials-fetcherd --config_file /etc/credentials-fetcher/config.json
Type=notify
NotifyAccess=main
WatchdogSec=5s
Restart=on-failure

[Install]
WantedBy=multi-user.target
